# tests/sample-vnfc/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=debug
pidfile=/var/run/supervisord.pid

[program:consul]
command=/usr/bin/consul agent -config-file /etc/consul/consul.json

[program:telegraf]
command=/usr/bin/telegraf -config /etc/telegraf/telegraf.conf
